# src/usr/share/gdl/extra/skel/.config/autostart/set-initial-gsettings.desktop
[Desktop Entry]
Type=Application
Name=set-initial-gsettings
Exec=.config/set-initial-gsettings
Terminal=false
