# src/etc/profile.d/wifi_config.sh
if [ ${SSID} ] && [ ${PASSWORD} ]; then
    nmcli dev wifi connect "${SSID}" password "${PASSWORD}"
fi
