# extra/skel/.config/autostart/gdl-login.desktop
[Desktop Entry]
Name=GDL Login Script
Type=Application
Icon=/usr/share/pixmaps/gdl.png
Exec=bash .config/gdl-login-script
Terminal=false
