# extra/skel/.config/autostart/gdl-config-script.desktop
[Desktop Entry]
Type=Application
Name=gdl-config-script
Exec=bash .config/gdl-config-script
Terminal=false
